# scripts/split/split.service
[Unit]
Description=Geth split.js auto distributor
After=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/your-chain
ExecStart=/usr/bin/geth attach %h/your-chain/node/validator1/geth.ipc --exec 'loadScript("scripts/split/split.js")'
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
